# .config/systemd/user/waybar.service.d/override.conf
[Service]
Environment=XDG_CURRENT_DESKTOP=Unity
ExecStart=
ExecStart=/usr/local/bin/waybar.sh