# .config/systemd/user/waybar.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/waybar -l debug
Restart=always
RestartSec=3