# .config/systemd/user/waybar.service.d/override.conf
[Service]
#Environment=XDG_CURRENT_DESKTOP=Unity
ExecStart=
ExecStart=/usr/bin/waybar --log-level debug
RestartSec=3