# .config/systemd/user/mako.service.d/override.conf
[Unit]
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Restart=always
ExecCondition=